# scripts/triggers/sut-provision.timer
[Unit]
Description=Run sut-provision periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
